# gerrit/files/gerrit.service
[Unit]
Description={{service_name}}

[Service]
ExecStart={{directory}}/bin/gerrit.sh start
ExecStop={{directory}}/bin/gerrit.sh stop
User={{user}}